# confs/supervisord.conf
[program:video]
numprocs=5
command=/opt/python/bin/python /data/video/scripts/task.py
process_name=%(program_name)s_%(process_num)02d
user=www-data
